# docker/app/config/supervisor.d/example.conf
[program:example]
command=php /var/www/artisan supervisor:example
autostart=true
autorestart=true